# target/linux/ramips/mt7621/base-files/etc/init.d/set-irq-affinity
#!/bin/sh /etc/rc.common

START=99

start() {
	if grep -q 'processor.*: 2' /proc/cpuinfo; then
		mask=4
	elif grep -q 'processor.*: 1' /proc/cpuinfo; then
		mask=2
	else
		return
	fi

	for irq in $(grep -E "mt|ra" /proc/interrupts | cut -d: -f1 | sed 's, *,,'); do
		echo "2" > "/proc/irq/$irq/smp_affinity"
	done

  /etc/init.d/wireless_startup.sh &
}
